# deploy/hdmi-kiosk-recovery.service
[Unit]
Description=Recover GNOME display manager if HDMI kiosk startup fails

[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl start gdm.service
