# deploy/hdmi-weston-smoke.service
[Unit]
Description=HDMI Weston EGL performance dashboard
After=hdmi-weston-kiosk.service
Requires=hdmi-weston-kiosk.service

[Service]
Type=simple
User=root
Environment=XDG_RUNTIME_DIR=/run/hdmi-kiosk
Environment=WAYLAND_DISPLAY=wayland-0
ExecStart=/home/p1ne4pp1e/screen_test/build/hdmi_egl_stress
Restart=on-failure
RestartSec=2

[Install]
WantedBy=multi-user.target
